# supervisor.conf
# =============================================================================
# XenonOS Laravel API - Supervisor Configuration for Queue Workers
# =============================================================================
# Place this file at: /etc/supervisor/conf.d/xenonos-worker.conf
# =============================================================================

[program:xenonos-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/xenonos-api/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/www/xenonos-api/storage/logs/worker.log
stopwaitsecs=3600
priority=10

# =============================================================================
# Laravel Reverb WebSocket Server (if using)
# =============================================================================

[program:xenonos-reverb]
command=php /var/www/xenonos-api/artisan reverb:start --host=127.0.0.1 --port=8080
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/xenonos-api/storage/logs/reverb.log
stopwaitsecs=10
priority=20

# =============================================================================
# Laravel Scheduler (alternative to cron)
# =============================================================================

[program:xenonos-scheduler]
command=php /var/www/xenonos-api/artisan schedule:work
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/xenonos-api/storage/logs/scheduler.log
stopwaitsecs=10
priority=30

# =============================================================================
# MANAGEMENT COMMANDS
# =============================================================================

# After creating this file:
# 1. Reload Supervisor configuration:
#    sudo supervisorctl reread
#    sudo supervisorctl update
#
# 2. Start all workers:
#    sudo supervisorctl start xenonos-worker:*
#
# 3. Start Reverb:
#    sudo supervisorctl start xenonos-reverb
#
# 4. Start Scheduler:
#    sudo supervisorctl start xenonos-scheduler
#
# 5. Check status:
#    sudo supervisorctl status
#
# 6. Restart all:
#    sudo supervisorctl restart all
#
# 7. Stop all:
#    sudo supervisorctl stop all
#
# 8. View logs:
#    tail -f /var/www/xenonos-api/storage/logs/worker.log
#    tail -f /var/www/xenonos-api/storage/logs/reverb.log
#    tail -f /var/www/xenonos-api/storage/logs/scheduler.log
